# password_notifier/cron.daily
#!/bin/bash

/opt/password_notifier/password_notifier.py > /dev/null 2>&1
